# 固件相关/原有BOOT分区的修订/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/sh
#communication 4 lcd4linux
::respawn:/bin/sh /etc/gserial.sh 
::ctrlaltdel:/bin/umount -a -r
